# etc/init.d/blk-availability
#!/sbin/openrc-run

description="Manages availability of block devices during shutdown"

depend() {
    keyword -softlevel
}

start() {
    : # This service has no start action.
}

stop() {
    ebegin "Deactivating LVM block devices"
    /usr/sbin/blkdeactivate -u -l wholevg -m disablequeueing -r wait
    eend $?
}
