# etc/init.d/lvm2-lvmpolld
#!/sbin/openrc-run

description="LVM2 poll daemon"
command="/usr/sbin/lvmpolld"
command_args="-t 60 -f"
supervisor=supervise-daemon

depend() {
	need dmeventd
	provide lvmpolld
}

start_pre() {
	checkpath --directory --mode 0700 /run/lvm
}
