# etc/init.d/lvm2-monitor
#!/sbin/openrc-run

description="Monitoring of LVM2 mirrors, snapshots etc."

depend() {
	need dmeventd
	before fsck localmount
	provide lvm-monitor
}

start() {
	ebegin "Activating LVM2 monitoring"
	lvm vgchange --monitor y
	eend $?
}

stop() {
	ebegin "Deactivating LVM2 monitoring"
	lvm vgchange --monitor n
	eend $?
}
